# deb/DEBIAN/postinst
#!/bin/bash

systemctl --system daemon-reload
systemctl enable ramlog-oe.service
echo "#####         Reboot to activate ramlog-oe         #####"
