# deb/DEBIAN/prerm
systemctl stop ramlog-oe
systemctl disable ramlog-oe
